# verysync-1.5.4-x5-5.0/bin/etc/linux-systemd/system/verysync-resume.service
[Unit]
Description=Restart Verysync after resume
Documentation=man:verysync(1)
After=sleep.target

[Service]
Type=oneshot
ExecStart=-/usr/bin/pkill -HUP -x verysync

[Install]
WantedBy=sleep.target
